# vagrant/zookeeper.conf
start on started toxiproxy
stop on stopping toxiproxy

exec ZK_DIR/bin/zkServer.sh start ZK_CONF
